# lims2/supervisor.env-server.conf
[program:env-server]
command=/bin/su -l -c "cd $HOME/env-server && exec node env-server.js" genee
