# lims2/supervisor.reserv-server.conf
[program:reserv-server]
command=/bin/su -l -c "cd $HOME/reserv-server && exec node index.js" genee
