# lims2/supervisor.epc-server.conf
[program:epc-server]
command=/bin/su -l -c "cd $HOME/epc-server && exec node index.js" genee
