# lims2/supervisor.node-lims2.conf
[program:node-lims2]
command=/bin/su -l -c "cd $HOME/node-lims2 && exec node index.js" genee
stopasgroup=true
killasgroup=true
