# lims2/supervisor.gdoor-server.conf
[program:gdoor-server]
command=/bin/su -l -c "cd $HOME/gdoor-server && exec node bin/index.js" genee
