# lims2/supervisor.xinetd.conf
[program:xinetd]
command=/usr/sbin/xinetd -dontfork -pidfile "/var/run/xinetd.pid" -stayalive