# lims2/supervisor.slapd.conf
[program:slapd]
command=/usr/sbin/slapd -d 0 -F /etc/ldap/slapd.d
