# lims2/supervisor.redis.conf
[program:cron]
command=/usr/bin/redis-server
