# lims2/supervisor.glogon-server.conf
[program:glogon-server]
command=/bin/su -l -c "cd $HOME/glogon-server && exec node bin/index.js" genee
