# lims2/supervisor.cacs-server.conf
[program:cacs-server]
command=/bin/su -l -c "cd $HOME/cacs-server && exec node bin/cacs-server.js" genee
