# lims2/supervisor.easyrun.conf
[program:easyrun]
command=/usr/local/bin/easyrun
