# lims2/supervisor.glogon-server.conf
[program:glogon-server]
command=/bin/su -l -c "cd $HOME/glogon-server && exec node glogon-server.js" genee
