# sphinxsearch/supervisor.sphinxsearch.conf
[program:sphinxsearch]
command=/usr/bin/searchd --nodetach --config /etc/sphinxsearch/sphinx.conf.sh
